# ansible/playbooks/3x-ui-metrics-firewall-host.service.j2
[Unit]
Description=Restrict host ports 9091/9100 (xray-exporter, node_exporter) to Alloy IP (sweden-node)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/3x-ui-metrics-firewall-host.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
